# sergio.service
[Unit]
Description=Sergio Bot - AI-powered Trello automation
After=network.target

[Service]
Type=simple
# Runs as deploy; uses sudo to execute Claude CLI as claudeuser.
# Requires: /etc/sudoers.d/sergio (created by npm run setup)
User=deploy
WorkingDirectory=%h/sergio
EnvironmentFile=%h/sergio/.env
ExecStart=/usr/bin/npx tsx src/index.ts
Restart=on-failure
RestartSec=10

[Install]
WantedBy=multi-user.target
